# launchd/com.rysk.gh-token-env.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  ログイン時に gh の Keychain トークンを GUI セッション環境変数 GH_TOKEN に注入する。
  Dock から起動される cmux の PR ポーリングを認証付き (5,000 req/h) にするのが目的。
  トークンは plist には保存されず、毎ログイン時に gh の Keychain から読む。
  `gh auth switch` 後は launchctl start com.rysk.gh-token-env で再注入するか再ログイン。

  セキュリティ上の許容事項 (単一ユーザーの Mac 前提):
    launchctl setenv はスコープをアプリ単位で分けられないため GH_TOKEN は GUI セッション
    全体へ露出し、GUI から起動した全アプリが参照できる。また setenv 実行中はトークンが
    argv に載り ps から数ミリ秒観測可能。いずれも launchctl の仕様上回避できないため、
    gh の OAuth トークンを最小スコープに保つ前提で受け入れる。

  導入:
    launchctl load ~/Library/LaunchAgents/com.rysk.gh-token-env.plist
  解除:
    launchctl unload ~/Library/LaunchAgents/com.rysk.gh-token-env.plist
    launchctl unsetenv GH_TOKEN
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.rysk.gh-token-env</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>t="$(gh auth token 2>/dev/null)"; [ -n "$t" ] &amp;&amp; launchctl setenv GH_TOKEN "$t"</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/rysk/.local/share/mise/shims:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
